# coolbot.service
[Unit]
Description=CoolBot

[Service]
ExecStart=/home/lmaurice/CoolBot/main.py

[Install]
WantedBy=multi-user.target
